# .config/autostart/remap.desktop
[Desktop Entry]
Encoding=UTF-8
Name=remap
Comment=Remaps keys
Exec=remap
Type=Application
StartupNotify=false
Terminal=false
Hidden=false

